# opencode/opencode.service
[Unit]
Description=OpenCode headless server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/bash -lc '%h/.opencode/bin/opencode serve --hostname 127.0.0.1 --port 4096'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
